# packer/files/overlays/laptop/etc/profile.d/gopath.sh
export GOROOT=/home/ivan/tools/go
export GOPATH=/home/ivan/go
